# systemd/reoptimize-macbook.service
[Unit]
Description=Reoptimize MacBook after kernel update or boot
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /usr/local/bin/reoptimize-macbook.sh
